# pkg/DEBIAN/prerm
#!/bin/bash
echo "Stopping Lich..."
systemctl stop lich
systemctl disable lich
